# investigate-PCIe-lockup/investigate-pcie-lockup.service
[Service]
Type=simple
ExecStart=/usr/local/bin/run-investigate-pcie-lockup.sh
